# apps/home/deploy/kittyhome.service
[Unit]
Description=KittyPaw Home
After=network.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{REMOTE_DIR}}
EnvironmentFile={{REMOTE_DIR}}/.env
ExecStart={{REMOTE_DIR}}/kittyhome
RuntimeDirectory=kittyhome
RuntimeDirectoryMode=0755
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
